# brewpy/crontab
  * *  *   *   *     python /home/tobe/src/rpi/brewpy/monitor.py >> /var/log/dht11.log
